# systemd/fake-hwclock.service
[Unit]
Description=Fake Hardware Clock
Documentation=man:fake-hwclock(8)
Wants=fake-hwclock-checkin.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/true
ExecStop=/usr/bin/fake-hwclock

[Install]
WantedBy=sysinit.target
