# systemd/fake-hwclock.timer
[Unit]
Description=Update fake hardware clock time
PartOf=fake-hwclock.service

[Timer]
OnUnitActiveSec=15m
