# systemd/fake-hwclock.service
[Unit]
Description=Fake Hardware Clock
Documentation=man:fake-hwclock(8)
DefaultDependencies=no
Conflicts=shutdown.target
Wants=fake-hwclock.timer

[Service]
Type=simple
ExecStart=/usr/bin/fake-hwclock
ExecStop=/usr/bin/fake-hwclock

[Install]
WantedBy=multi-user.target
