# systemd/fake-hwclock-checkin.timer
[Unit]
Description=Update fake hardware clock time
PartOf=fake-hwclock.service fake-hwclock-checkin.service

[Timer]
OnUnitActiveSec=15m
